# templates/systemd/transmission-rpc-proxy/inside.socket.j2
[Unit]
Description=Transmission RPC inner proxy

[Socket]
ListenStream=%t/transmission-rpc.socket
SocketGroup={{ transmission.rpc.socket.group }}
SocketMode=0060

[Install]
WantedBy=sockets.target
